# dot_config/systemd/user/poweralertd.service.d/override.conf
[Unit]
After=graphical-session.target

[Service]
Type=exec
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "niri:wlroots:sway:Wayfire:labwc:Hyprland" ""
Slice=app-graphical.slice
